# config/service/sv/rivercarro/run
#!/bin/sh

exec 2>&1
exec rivercarro -view-padding 4 -outer-padding 0


