# config/service/sv/yambar-light/run
#!/bin/sh

exec 2>&1
exec yambar -c "$HOME"/.config/yambar/config-river-light.yml

