# config/service/sv/ssh-agent/run
#!/bin/sh

exec 2>&1
exec pgrep -f "^ssh-agent" && echo "ssh-agent already running"
exec ssh-agent -D -a "$SSH_AUTH_SOCK"

