# system/etc/sv/snooze-fstrim/run
#!/bin/sh

# Fstrim once a week
exec 2>&1
exec snooze -w0 -s 6d -t ./timefile sh -c 'fstrim -v / && touch timefile'
