# config/service/foot/run
#!/bin/sh

exec /usr/bin/foot --server 2>&1

