# config/service/sv/snooze-trash-empty/run
#!/bin/sh

# Purge trashed files older than 1 day
exec 2>&1
exec snooze -v -H10 -s 1d -t ./timefile sh -c 'trash-empty 1 && touch timefile'

